# supervisord.conf
[program:yellowent]
command=node app.js
directory=/opt/box