# yellowtent.conf
; Yellowtent supervisor config file.

[program:box]
command=node app.js
directory=/opt/box
