# supervisor/conf.d/box.conf
[program:box]
command=node app.js
autostart=false
autorestart=true
directory=%(ENV_NGINX_ROOT)s