# supervisor/conf.d/nginx.conf
[program:nginx]
command=nginx -c nginx.conf -p %(ENV_NGINX_ROOT)s/nginx/
autostart=true
autorestart=true