# supervisor/conf.d/yellowtent.conf
[program:box]
command=node app.js
autostart=false
autorestart=true