# setup/container/systemd/cloudron.target
[Unit]
Description=Cloudron Smart Cloud
Documentation=https://cloudron.io/documentation.html
StopWhenUnneeded=true
Requires=box.service janitor.timer docker_janitor.timer
After=box.service janitor.timer docker_janitor.timer
# AllowIsolate=yes

[Install]
WantedBy=multi-user.target
